# com.stockmonitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.stockmonitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/yiming/Documents/Claude/Projects/Etrade stock monitor/.venv/bin/python</string>
    <string>-m</string>
    <string>src.main</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/yiming/Documents/Claude/Projects/Etrade stock monitor</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/yiming/Documents/Claude/Projects/Etrade stock monitor/logs/launchd_stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/yiming/Documents/Claude/Projects/Etrade stock monitor/logs/launchd_stderr.log</string>

  <key>ThrottleInterval</key>
  <integer>30</integer>

</dict>
</plist>
